# rc.local
#!/bin/bash
/etc/mine/start.sh &
exit 0
